# tools/airship-in-a-pod/libvirt/assets/usr/lib/systemd/system/libvirtd.service.d/10-logs.conf
[Service]
StandardOutput=tty
#FailureAction=poweroff
